# ppetrov/openvpn/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/openvpn/
drwxr-x---	root/nobody	etc/openvpn/certs/
drwxr-x---	root/nobody	etc/openvpn/keys/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/openvpn
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/openvpn-msg.h
-rw-r--r--	root/root	usr/include/openvpn-plugin.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/openvpn/
drwxr-xr-x	root/root	usr/lib/openvpn/plugins/
-rwxr-xr-x	root/root	usr/lib/openvpn/plugins/openvpn-plugin-down-root.la
-rwxr-xr-x	root/root	usr/lib/openvpn/plugins/openvpn-plugin-down-root.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/openvpn
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/openvpn-examples.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/openvpn.8.gz
drwxr-xr-x	root/root	usr/share/openvpn/
-rw-r--r--	root/root	usr/share/openvpn/openvpn.conf
drwxrwxr-x	root/root	usr/share/openvpn/sample-config-files/
-rw-rw-r--	root/root	usr/share/openvpn/sample-config-files/README
-rw-rw-r--	root/root	usr/share/openvpn/sample-config-files/client.conf
-rwxrwxr-x	root/root	usr/share/openvpn/sample-config-files/firewall.sh
-rwxrwxr-x	root/root	usr/share/openvpn/sample-config-files/home.up
-rw-rw-r--	root/root	usr/share/openvpn/sample-config-files/loopback-client
-rw-rw-r--	root/root	usr/share/openvpn/sample-config-files/loopback-server
-rwxrwxr-x	root/root	usr/share/openvpn/sample-config-files/office.up
-rwxrwxr-x	root/root	usr/share/openvpn/sample-config-files/openvpn-shutdown.sh
-rwxrwxr-x	root/root	usr/share/openvpn/sample-config-files/openvpn-startup.sh
-rw-rw-r--	root/root	usr/share/openvpn/sample-config-files/server.conf
-rw-rw-r--	root/root	usr/share/openvpn/sample-config-files/tls-home.conf
-rw-rw-r--	root/root	usr/share/openvpn/sample-config-files/tls-office.conf
drwxrwxr-x	root/root	usr/share/openvpn/sample-keys/
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/README
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/ca.crt
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/ca.key
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/client-ec.crt
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/client-ec.key
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/client-pass.key
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/client.crt
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/client.key
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/client.p12
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/dh2048.pem
-rwxrwxr-x	root/root	usr/share/openvpn/sample-keys/gen-sample-keys.sh
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/openssl.cnf
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/server-ec.crt
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/server-ec.key
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/server.crt
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/server.key
-rw-rw-r--	root/root	usr/share/openvpn/sample-keys/ta.key
drwxrwxr-x	root/root	usr/share/openvpn/sample-scripts/
-rwxrwxr-x	root/root	usr/share/openvpn/sample-scripts/auth-pam.pl
-rwxrwxr-x	root/root	usr/share/openvpn/sample-scripts/bridge-start
-rwxrwxr-x	root/root	usr/share/openvpn/sample-scripts/bridge-stop
-rwxrwxr-x	root/root	usr/share/openvpn/sample-scripts/totpauth.py
-rwxrwxr-x	root/root	usr/share/openvpn/sample-scripts/ucn.pl
-rwxrwxr-x	root/root	usr/share/openvpn/sample-scripts/verify-cn
